# vagrant/systemd/sources/nginxMemoryLimit.slice
[Unit]
Description=Slice that limits memory for nginx services

[Slice]
MemoryMax=50M
